# files/etc/uci-defaults/50-tessel-hostname
#!/bin/sh

MACADDR=$(ifconfig eth0 | awk '/HWaddr/ { print $5 }' | tr -d :)

uci batch <<EOF
set system.@system[0].hostname=Tessel-${MACADDR}
commit system
EOF
